# musicbox.service
[Unit]
Description=Music Box

[Service]
Type=simple
ExecStart=/home/pi/Documents/Music-Box/main.py