# musicbox.service
[Unit]
Description=Music Box

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u main.py
WorkingDirectory=/home/pi/Documents/Music-Box